# systemd/windrose-bot.service
# windrose-bot.service -- Windrose Telegram control bot
# ADR-006 (python-telegram-bot v22), ADR-007 (access control)
#
# BOT_USER_PLACEHOLDER is stamped with the real username by install_service.sh.

[Unit]
Description=Windrose Telegram Bot
Wants=network-online.target
After=syslog.target network-online.target windrose.service

[Service]
Type=simple
User=BOT_USER_PLACEHOLDER
WorkingDirectory=/home/BOT_USER_PLACEHOLDER/windrose-telegram-bot

ExecStart=/home/BOT_USER_PLACEHOLDER/windrose-telegram-bot/venv/bin/python \
          /home/BOT_USER_PLACEHOLDER/windrose-telegram-bot/bot.py

Restart=on-failure
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=windrose-bot

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=tmpfs
BindPaths=/home/BOT_USER_PLACEHOLDER/windrose-telegram-bot \
          /home/BOT_USER_PLACEHOLDER/log \
          /home/BOT_USER_PLACEHOLDER/scripts \
          /home/BOT_USER_PLACEHOLDER/windrose-backups

[Install]
WantedBy=multi-user.target
